# root/etc/s6-overlay/s6-rc.d/init-etherpad-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
#
# On first run, seed /config/settings.json from the bundled template,
# redirect /app/etherpad-lite/var into /config so installed_plugins.json
# and runtime state persist, pre-seed installed_plugins.json (so
# Etherpad's startup checkForMigration doesn't shell out to `pnpm ls`
# from a directory abc can't write to), then optionally install plugins
# listed in /config/plugins.txt.

APP_DIR="${EP_DIR:-/app/etherpad-lite}"
SETTINGS="/config/settings.json"
TEMPLATE="${APP_DIR}/settings.json.docker"
[ -f "${TEMPLATE}" ] || TEMPLATE="${APP_DIR}/settings.json.template"

# 1. Seed settings.json on first launch.  The template uses Etherpad's
#    native ${VAR:default} substitution syntax, so DB selection / IP /
#    port etc. are driven by env vars set in svc-etherpad/run, NOT by
#    rewriting this file.  An admin editing /config/settings.json to a
#    literal value (e.g. "dbType": "postgres") will override the env
#    default cleanly.
if [ ! -f "${SETTINGS}" ]; then
    echo "[etherpad] bootstrapping ${SETTINGS} from ${TEMPLATE}"
    cp "${TEMPLATE}" "${SETTINGS}"
fi

# 2. Redirect the app's var/ into /config so installed_plugins.json,
#    sessionkey, APIKEY etc. survive container recreation.  The Etherpad
#    plugin installer uses __dirname to locate its registry, so we
#    cannot point it at /config via env vars alone — a symlink under
#    the app tree is the only thing the runtime honours.
mkdir -p /config/var
if [ -L "${APP_DIR}/var" ]; then
    rm -f "${APP_DIR}/var"
elif [ -e "${APP_DIR}/var" ]; then
    # Carry forward anything the build left in var/ (e.g. APIKEY.txt) the
    # first time around, then replace with the symlink.
    cp -a "${APP_DIR}/var/." /config/var/ 2>/dev/null || true
    rm -rf "${APP_DIR}/var"
fi
ln -s /config/var "${APP_DIR}/var"

# 3. Pre-seed installed_plugins.json so Etherpad's checkForMigration
#    skips its `pnpm ls --long --json --depth=0 --no-production` call
#    on boot.  pnpm needs to write a tmp file under ${APP_DIR}, which
#    is root-owned by the image build, and abc cannot write there.
#    The seed only writes the bundled core plugin; the admin UI adds
#    real entries when users install plugins at runtime.
INSTALLED="/config/var/installed_plugins.json"
if [ ! -f "${INSTALLED}" ]; then
    VERSION=$(node -p "require('${APP_DIR}/src/package.json').version" 2>/dev/null \
              || node -p "require('${APP_DIR}/package.json').version" 2>/dev/null \
              || echo "0.0.0")
    cat > "${INSTALLED}" <<EOF
{"plugins":[{"name":"ep_etherpad-lite","version":"${VERSION}"}]}
EOF
fi

# 3a. Etherpad writes APIKEY.txt and SESSIONKEY.txt directly under
#     settings.root (= ${APP_DIR}) on first boot. That directory is
#     root-owned by the image build, so abc can't write — API requests
#     fail with "no or wrong API Key" and admin pad-session resets
#     break. Symlink both keyfiles into /config so they're abc-writable
#     AND persist across container recreation.
#
#     IMPORTANT: do NOT pre-touch the /config target empty. Etherpad's
#     generate-on-missing logic only fires when readFileSync throws —
#     an empty file reads as the empty string, so Etherpad treats it
#     as a valid (but blank) key and every API request fails with
#     "no or wrong API Key". Leave the target absent so the first
#     write through the symlink populates it correctly.
for keyfile in APIKEY.txt SESSIONKEY.txt; do
    target="/config/${keyfile}"
    link="${APP_DIR}/${keyfile}"
    if [ ! -L "${link}" ]; then
        # Carry forward a non-empty build-time file (rare; the upstream
        # image normally doesn't ship one) before we drop the symlink in.
        if [ -s "${link}" ] && [ ! -s "${target}" ]; then
            cp "${link}" "${target}"
            chown abc:abc "${target}"
            chmod 0640 "${target}"
        fi
        rm -f "${link}"
        ln -s "${target}" "${link}"
    fi
done

# 4. Install any user-requested plugins.  /config/plugins.txt holds one
#    plugin name per line; lines starting with '#' or blank are skipped.
#    pnpm wants HOME writable for its store, and s6-setuidgid leaves
#    HOME at /root, so pass an abc-writable HOME via env.
if [ -f /config/plugins.txt ]; then
    cd "${APP_DIR}" || exit 1
    while IFS= read -r plugin; do
        plugin="${plugin%%#*}"
        plugin="$(echo "${plugin}" | xargs)"
        [ -z "${plugin}" ] && continue
        if [ ! -d "${APP_DIR}/src/node_modules/${plugin}" ]; then
            echo "[etherpad] installing plugin ${plugin}"
            s6-setuidgid abc env HOME=/config \
                pnpm run plugins i "${plugin}" || \
                echo "[etherpad] WARNING: failed to install ${plugin}"
        fi
    done < /config/plugins.txt
fi

# 5. Final ownership pass.  /config/* and the plugin-writable dirs under
#    /app/etherpad-lite need to be abc-owned so Etherpad can write to
#    them.  plugin_packages and node_modules are root-owned by the build
#    but the admin UI installs new plugins into them at runtime.
lsiown -R abc:abc /config
for dir in \
    "${APP_DIR}/src/node_modules" \
    "${APP_DIR}/src/plugin_packages"; do
    if [ -d "${dir}" ]; then
        lsiown -R abc:abc "${dir}"
    fi
done

# Force success — see prior commit notes; the optional-dir loop above
# can leave a non-zero exit status if the last directory tested doesn't
# exist on a fresh install, which s6-rc would read as init failure.
exit 0
