# root/etc/s6-overlay/s6-rc.d/svc-etherpad/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
#
# Launch Etherpad directly via node + tsx, bypassing the pnpm runtime
# wrapper.  pnpm 11.x's `run` invocation triggers a "deps status check"
# that writes a tmp file under the app root, which is root-owned by
# the build — abc cannot write there and the service dies before
# Etherpad ever boots.  node is on PATH (apk-installed in the image),
# tsx is in src/node_modules from the production install, so there is
# no behavioural difference vs `pnpm run prod` other than the missing
# verification step.

APP_DIR="${EP_DIR:-/app/etherpad-lite}"

# HOME for anything that needs it (Etherpad itself doesn't, but some
# plugins or downstream tools may).
export HOME=/config

# Defaults that override the dev-only template values for sqlite at a
# writable persistent path.  Users who set DB_TYPE / DB_FILENAME via
# `docker run -e` win — the with-contenv wrapper exports those into our
# environment, and `:-` only fills the default when unset/empty.
export DB_TYPE="${DB_TYPE:-sqlite}"
export DB_FILENAME="${DB_FILENAME:-/config/etherpad.db}"

# Default to apikey auth so the HTTP API works out of the box (Etherpad
# only auto-generates APIKEY.txt when authenticationMethod === 'apikey';
# the template default is 'sso' which never writes the file and leaves
# every /api/1/* request returning "no or wrong API Key"). Users who
# actually want SSO can set AUTHENTICATION_METHOD=sso via env.
export AUTHENTICATION_METHOD="${AUTHENTICATION_METHOD:-apikey}"

cd "${APP_DIR}/src" || exit 1

# s6-notifyoncheck was here originally (matching the thelounge template
# we cloned from), but it exits "s6-supervise not running" before node
# can start in our setup, killing svc-etherpad in a tight restart loop.
# Drop it — s6 treats process-alive as readiness, which is what we want
# for a long-running Etherpad daemon. Healthchecking happens at the
# docker layer (HEALTHCHECK below in Dockerfile or compose-side).
exec \
    s6-setuidgid abc \
    node --import tsx node/server.ts --settings /config/settings.json
