# root/etc/s6-overlay/s6-rc.d/svc-etherpad/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
#
# Launch Etherpad directly via node + tsx, bypassing the pnpm runtime
# wrapper.  pnpm 11.x's `run` invocation triggers a "deps status check"
# that writes a tmp file under the app root, which is root-owned by
# the build — abc cannot write there and the service dies before
# Etherpad ever boots.  node is on PATH (apk-installed in the image),
# tsx is in src/node_modules from the production install, so there is
# no behavioural difference vs `pnpm run prod` other than the missing
# verification step.

APP_DIR="${EP_DIR:-/app/etherpad-lite}"

# HOME for anything that needs it (Etherpad itself doesn't, but some
# plugins or downstream tools may).
export HOME=/config

# Defaults that override the dev-only template values for sqlite at a
# writable persistent path.  Users who set DB_TYPE / DB_FILENAME via
# `docker run -e` win — the with-contenv wrapper exports those into our
# environment, and `:-` only fills the default when unset/empty.
export DB_TYPE="${DB_TYPE:-sqlite}"
export DB_FILENAME="${DB_FILENAME:-/config/etherpad.db}"

cd "${APP_DIR}/src" || exit 1

exec \
    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost 9001" \
    s6-setuidgid abc \
    node --import tsx node/server.ts --settings /config/settings.json
