# root/etc/s6-overlay/s6-rc.d/svc-etherpad/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# pnpm reads ~/.npmrc and writes its store under $HOME/.local/share/pnpm.
# s6-setuidgid only changes uid/gid; HOME stays at /root, which abc
# can't write to.  Point HOME at /config (abc-owned by init script) so
# pnpm has somewhere to land.
export HOME=/config

cd "${EP_DIR:-/app/etherpad-lite}" || exit 1

exec \
    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost 9001" \
    s6-setuidgid abc \
    env HOME=/config \
    pnpm run prod -- --settings /config/settings.json
